# home/etc/profile.d/editor.sh
# vim: nu: tabstop=8: softtabstop=8: shiftwidth=8

EDITOR=/usr/local/bin/vim
