# home/etc/profile.d/git-prompt-settings.sh
# vim: nu: tabstop=8: softtabstop=8: shiftwidth=8

# Settings for git-prompt.sh
GIT_PS1_SHOWSTATUSSTRING="on"
GIT_PS1_SHOWUPSTREAM="git verbose"
GIT_PS1_SHOWDIRTYSTATE="on"
GIT_PS1_SHOWSTASHSTATE="on"
GIT_PS1_SHOWCOLORHINTS="on"
